# freebsd10/usr/local/etc/rc.d/cloudsshkey
#!/bin/sh
#
# cloudsshkey:
#  Allows setting ssh public keys
#  via LeaseWeb Self Service Center

# PROVIDE: cloudsshkey
# REQUIRE: LOGIN NETWORKING cleanvar
# KEYWORD: shutdown

. /etc/rc.subr
name="cloudsshkey"
command="/usr/local/sbin/cloudsshkey"
rcvar=cloudsshkey_enable
load_rc_config "${name}"
run_rc_command "$1"

